# system_files/shared/usr/lib/systemd/user/rpm-ostree-advisory-notifier.timer
[Unit]
Description=Daily User RPM OSTree Security Advisory Notifier Prompt

[Timer]
OnBootSec=2h
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target