# system_files/shared/usr/lib/systemd/user/rpm-ostree-advisory-notifier.service
[Unit]
Description=RPM OSTree Security Advisory Notifier

[Service]
Type=simple
ExecStart=/usr/libexec/rpm-ostree-advisory-notifier